# ansible/templates/asterisk-worker.service.j2
[Unit]
Description=Asterisk Node.js Worker - AMI/ARI Controller
After=network.target asterisk.service redis6.service
Wants=asterisk.service redis6.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/asterisk-worker
ExecStart=/usr/bin/node --max-old-space-size=512 /opt/asterisk-worker/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asterisk-worker

# Environment variables
Environment="NODE_ENV=production"
Environment="PORT={{ nodejs_worker_port }}"
Environment="AMI_HOST=localhost"
Environment="AMI_PORT={{ ami_port }}"
Environment="AMI_USERNAME={{ ami_username }}"
Environment="AMI_PASSWORD={{ ami_password }}"
Environment="ARI_HOST=localhost"
Environment="ARI_PORT={{ ari_port }}"
Environment="ARI_USERNAME={{ ari_username }}"
Environment="ARI_PASSWORD={{ ari_password }}"
Environment="AWS_REGION={{ aws_region }}"
Environment="S3_AUDIO_BUCKET={{ s3_audio_bucket }}"

[Install]
WantedBy=multi-user.target
